# example.crontab
0 4 * * * /full/path/to/dump-ac-status/run.sh >> /full/path/to/dump-ac-status.log 2>>&1
